# apps/runner/packaging/deb/DEBIAN/postinst
#!/bin/bash
set -e

# Create necessary directories
mkdir -p /var/lib/nightona/runner
mkdir -p /var/log/nightona
mkdir -p /etc/nightona

# Create default environment file if it doesn't exist
if [ ! -f /etc/nightona/runner.env ]; then
    cat > /etc/nightona/runner.env <<EOF
# Nightona Runner Environment Variables
# Add your configuration here, for example:
# NIGHTONA_API_URL=https://api.daytona.io
# NIGHTONA_RUNNER_TOKEN=your-token-here
EOF
    chmod 600 /etc/nightona/runner.env
fi

# Reload systemd daemon
systemctl daemon-reload

echo "Nightona Runner installed successfully"
echo ""
echo "Next steps:"
echo "1. Configure environment variables in /etc/nightona/runner.env"
echo "2. Enable and start the service: systemctl enable --now nightona-runner"
echo ""
echo "To check service status: systemctl status nightona-runner"

exit 0
